# files/etc/s6-overlay/s6-rc.d/user-setup/up
#!/command/with-contenv bash

exec /app/user-setup.sh